# packer/files/boincme-config/boincme-config.service
[Service]
ExecStart=/boincme-config.sh
[Install]
WantedBy=default.target