# packer/boinc-config/boinc-config.service
[Service]
ExecStart=/boinc-config.sh
[Install]
WantedBy=default.target